# avnet-iot-1.2/services/attinyupdate.service
[Unit]
Description=Controls ATTINY flash update 

[Service]
Type=simple
ExecStart=/opt/avnet-iot/iotservices/attiny_flash

[Install]
WantedBy=default.target
